# images/ochostats/resources/supervisor/ochostats_standalone.conf
# No Ochostats in standalone mode
#[program:ochostats]
#command=/usr/bin/python /opt/ochostats/pod/pod.py

[program:statsd]
command=/usr/bin/nodejs /opt/statsd/stats.js /etc/statsd/localConfig.js

[program:carbon]
command=/usr/bin/python /usr/bin/carbon-cache.py --config=/etc/carbon/carbon.conf --pidfile=/var/run/carbon-cache.pid --logdir=/var/log/carbon/ --nodaemon start

[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND

[program:grafana]
directory=/usr/share/grafana
command=/usr/sbin/grafana-server --pidfile=/var/run/grafana-server.pid --config=/etc/grafana/grafana.ini cfg:default.paths.data=/var/lib/grafana cfg:default.paths.logs=/var/log/grafana

[program:postgres]
command=sudo -u postgres /usr/lib/postgresql/9.3/bin/postgres -D /var/lib/postgresql/9.3/main -c config_file=/etc/postgresql/9.3/main/postgresql.conf
